# puppet/private_modules/cikl/templates/cikl-dns-worker-upstart.conf.erb
description "Cikl DNS Worker"

start on vagrant-mounted MOUNTPOINT=/vagrant

stop on runlevel [06]

respawn

setuid <%= scope.lookupvar('cikl::worker_user') %>
setgid <%= scope.lookupvar('cikl::worker_group') %>

console log
chdir /vagrant/cikl-worker
exec bin/dns_worker.rb /etc/cikl-dns-worker.yaml
