# puppet/private_modules/cikl/templates/cikl-dns-worker-upstart.conf.erb
description "Cikl DNS Worker"

start on started networking
stop on runlevel [06]

respawn

#setuid cikl
#setgid cikl

console log
chdir /vagrant/cikl-worker
exec bin/dns_worker.rb /etc/cikl-dns-worker.yaml
