# puppet/private_modules/profile/templates/worker/cikl-dns-worker-upstart.conf.erb
description "Cikl DNS Worker"

start on (started networking <%= "and vagrant-mounted MOUNTPOINT=/vagrant" if (@environment == 'development') %>)

stop on runlevel [06]

respawn

setuid <%= @user %>
setgid <%= @group %>

console log
chdir <%= @local_path %>
exec bin/dns_worker.rb /etc/cikl-dns-worker.yaml
