# puppet/private_modules/profile/templates/worker/cikl-dns-worker-upstart.conf.erb
description "Cikl DNS Worker"

start on (started networking <%= "and vagrant-mounted MOUNTPOINT=/vagrant" if (@environment == 'development') %>)

stop on runlevel [06]

respawn

setuid <%= @user %>
setgid <%= @group %>

console log
env BUNDLE_GEMFILE=<%= @local_path %>/Gemfile
chdir <%= @server_run_path %>
exec <%= @local_path %>/bin/dns_worker.rb <%= @server_config %>
