# puppet/private_modules/profile/templates/api/cikl-api-pre-upstart.conf.erb
description "Cikl API Pre-start"

# Creates required directories before starting up our service.

start on starting cikl-api

task

script
  mkdir -p <%= @server_run_path %>
  chown <%= @user %>:<%= @group %> <%= @server_run_path %>

  mkdir -p <%= @server_log_path %>
  chown <%= @user %>:<%= @group %> <%= @server_log_path %>
end script
