# puppet/private_modules/profile/templates/api/cikl-api-upstart.conf.erb
description "Cikl API"

start on (started networking <%= "and vagrant-mounted MOUNTPOINT=/vagrant" if (@environment == 'development') %>)

stop on runlevel [06]

setuid <%= @user %>
setgid <%= @group %>

console log
script
  export APP_ROOT=<%= @local_path %>
  export BUNDLE_GEMFILE=<%= @local_path %>/Gemfile
  chdir <%= @server_run_path %>
  bundle exec puma -e <%= @environment %> -C <%= @server_config %>
end script

